# faceapp.conf
[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as in [*_http_server] if set
;password=123                ; should be same as in [*_http_server] if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

; The sample program section below shows all possible program subsection values.
; Create one or more 'real' program: sections to be able to control them under
; supervisor.


[program:face_core]
command=/home/firebird/pythonProject/flasky/venv/bin/gunicorn -c gunicorn.py wsgi:app    ; supervisor启动命令
directory=/home/firebird/pythonProject/flasky/faceService                               ; 项目的文件夹路径
startsecs=0                                                                             ; 启动时间
stopwaitsecs=0                                                                          ; 终止等待时间
autostart=false                                                                         ; 是否自动启动
autorestart=false                                                                       ; 是否自动重启
stdout_logfile=/home/firebird/pythonProject/flasky/faceService/logs/gunicorn.access.log                           ; log 日志
stderr_logfile=/home/firebird/pythonProject/flasky/faceService/logs/gunicorn.error.log                          ; 错误日志
